# edi/plugins/playbooks/debian/base_system_cleanup/roles/openssh_server_keys/files/openssh-server-key-generator.service
[Unit]
Description=Generate host keys for openssh-server.
Before=ssh.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A -v

[Install]
RequiredBy=multi-user.target

